# src/inittab
::ctrlaltdel:/sbin/reboot
::shutdown:/bin/umount -a -r
::restart:/sbin/init
ttyS0::sysinit:/etc/init.d/rcS
ttyS0::wait:/usr/local/bin/scopeApp.sh
ttyS0::respawn:/sbin/getty -L ttyS0 115200
ttyS0::once:/bin/mv /usr/local/nv/nvkey20.dat /usr/local/nv/nvkey20.crash
